# docker/s6/php-fpm/run
#!/bin/sh

/usr/local/sbin/php-fpm --pid /var/run/php-fpm.pid -F
